# conf/systemd/nginx.service.d/override.conf
[Unit]
OnFailure=email-failures@nginx.service

[Service]
# let nginx decide on that
LimitNOFILE=65535
Restart=always
