# conf/systemd/nginx.service.d/override.conf
[Unit]
OnFailure=email-failures@nginx-%l.service